# .docker/etc/supervisord.d/vault.ini
[program:vault]
command=/usr/bin/vault agent -config=/etc/vault-agent/agent-config.hcl -log-level=debug 
numprocs=1                    ; number of processes copies to start (def 1)
umask=022                     ; umask for process (default None)
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; retstart at unexpected quit (default: true)
user=root                     ; (default is current user, required if root)